# systemd/auto-upgrade.service
[Unit]
Description=Auto-upgrade system packages
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/intelluxe/scripts/auto-upgrade.sh
User=root
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=yes
SyslogIdentifier=auto-upgrade
PrivateTmp=yes
