# systemd/bootstrap.service
# bootstrap.service version 1.0.0
# Author: Justin Sue (Galdaer)
# Last Revised: 2025-05-30
# Purpose: Systemd unit for clinic + WireGuard stack setup.
# Usage: Brings up the main stack, used at boot or after updates.

[Unit]
Description=clinic + WireGuard Setup
After=network-online.target docker.service
Wants=docker.service
# Don't block graphical.target startup
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/opt/intelluxe/scripts/bootstrap.sh --log-file /var/log/bootstrap.log
# Remove restart for oneshot - use timer instead
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bootstrap
Environment=DEBUG=false
Environment=DRY_RUN=false
Environment=COLOR=true

# (Optional) Give network tools more time
TimeoutStartSec=120

[Install]
# Start after graphical session is established
WantedBy=graphical-session.target
Also=bootstrap.timer
