# systemd/clinic-teardown.service
# clinic-teardown.service version 1.0.0
# Author: Justin Sue (Galdaer)
# Last Revised: 2025-05-30
# Purpose: Systemd unit for clinic stack teardown utility.
# Usage: Runs clinic-teardown.sh for full stack cleanup.

[Unit]
Description=Run clinic-teardown.sh cleanup utility
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/clinic-teardown.sh --all --force
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
