# systemd/clinic-wg.service
# clinic-wg.service version 1.0.0
# Author: Justin Sue (Galdaer)
# Last Revised: 2025-05-30
# Purpose: Systemd unit to bring up the WireGuard interface within the SHaN namespace.
# Usage: Runs clinic-wg-setup.sh which executes wg-quick inside the namespace.

[Unit]
Description=Bring up WireGuard interface inside SHaN namespace
After=clinic-netns.service
Requires=clinic-netns.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/clinic-wg-setup.sh

[Install]
WantedBy=multi-user.target
