# systemd/clinic-diagnostics.service
# clinic-diagnostics.service version 1.0.0
# Author: Justin Sue (Galdaer)
# Last Revised: 2025-05-30
# Purpose: Systemd unit for clinic Diagnostics Health Check.
# Usage: Used with clinic-diagnostics.timer for scheduled health checks

[Unit]
Description=clinic Diagnostics Health Check
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/clinic-diagnostics.sh --log-file /var/log/clinic-diagnostics.log --deep-check
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clinic-diagnostics
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
