# systemd/diagnostics.service
# diagnostics.service version 1.0.0
# Author: Justin Sue (Galdaer)
# Last Revised: 2025-05-30
# Purpose: Systemd unit for clinic Diagnostics Health Check.
# Usage: Used with diagnostics.timer for scheduled health checks

[Unit]
Description=clinic Diagnostics Health Check
After=network-online.target docker.service graphical.target
Wants=docker.service
# Don't block boot sequence
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/opt/intelluxe/scripts/diagnostics.sh --log-file /var/log/diagnostics.log --deep-check
# Use timer for restarts instead of service restart
RemainAfterExit=no
StandardOutput=journal
StandardError=journal
SyslogIdentifier=diagnostics
TimeoutStartSec=120

[Install]
# Only start via timer, not at boot
WantedBy=
