# systemd/duckdns-update.timer
# duckdns-update.timer version 1.0.0
# Author: Justin Sue (Galdaer)
# Last Revised: 2025-05-30
# Purpose: Systemd timer for periodic DuckDNS IP updates.
# Usage: Triggers duckdns-update.service every 5 minutes.

[Unit]
Description=Run DuckDNS updater every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=1day
Persistent=true
Unit=duckdns-update.service

[Install]
WantedBy=timers.target
