# systemd/resource-pusher.timer
# resource-pusher.timer version 1.0.0
# Author: Justin Sue (Galdaer)
# Last Revised: 2025-05-30
# Purpose: Systemd timer for periodic host resource metric pushes.
# Usage: Triggers resource-pusher.service every 10 minutes.

[Unit]
Description=Run resource-pusher every 10 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
Unit=intelluxe-resource-pusher.service

[Install]
WantedBy=timers.target
